# pi-setup/mk3-screen-daemon.service
[Unit]
Description=MK3 Screen Mirroring Daemon
# Start early — mirrors whatever is on /dev/fb0 (boot console, then Mixxx)
After=local-fs.target
Wants=local-fs.target

[Service]
Type=simple
# Run as root for /dev/fb0 access during early boot.
# The udev rule grants USB access to the audio group.
ExecStart=/usr/local/bin/mk3-screen-daemon --resolution 960x544 --fps 15 --no-partial
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
