# pi-setup/mk3-screen-daemon.service
[Unit]
Description=MK3 Screen Mirroring Daemon
After=mixxx.service
Requires=mixxx.service

[Service]
Type=simple
User=pi
ExecStart=/usr/local/bin/mk3-screen-daemon --resolution 960x544 --fps 30 --no-partial
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
